# nftblocklist.cron.daily
#!/bin/bash

/usr/sbin/nftblocklist -v -r
